# ops/systemd/lumii-backend.service.d/40-public-api-https.conf
[Service]
Environment="LUMII_BACKEND_HOST=127.0.0.1"
Environment="LUMII_PUBLIC_API_BASE_URL=https://api.lumiiapp.cn"
Environment="LUMII_PUBLIC_API_PROBE_TIMEOUT_MS=6000"
Environment="LUMII_PUBLIC_API_PROBE_CONNECT_ADDRESS=127.0.0.1"
