# ops/systemd/lumii-backend.service.d/45-media-cdn.example.conf
[Service]
# Public media URLs returned to the App.
Environment="PET_AVATAR_PUBLIC_BASE_URL=https://media.lumiiapp.cn"

# Independent media-domain probe shown by GET /admin/system/health.
Environment="MEDIA_PUBLIC_PROBE_BASE_URL=https://media.lumiiapp.cn"
Environment="MEDIA_PUBLIC_PROBE_TIMEOUT_MS=6000"

# Public media is browser-cacheable for at most 60 seconds and shared/CDN
# cacheable for at most 300 seconds until an active CDN purge workflow exists.
Environment="COS_PROXY_BROWSER_CACHE_SECONDS=60"
Environment="COS_PROXY_CACHE_SECONDS=300"
