# ops/systemd/lumii-backend.service.d/50-sqlite-state.conf
[Service]
Environment="NODE_ENV=production"
Environment="LUMII_BACKEND_STATE_PATH=/home/ubuntu/lumii-data/lumii-backend-state.json"
Environment="LUMII_STATE_STORAGE_DRIVER=sqlite"
Environment="LUMII_STATE_SQLITE_PATH=/home/ubuntu/lumii-data/lumii-state.sqlite"
Environment="LUMII_STATE_SQLITE_BUSY_TIMEOUT_MS=5000"
Environment="LUMII_STATE_SQLITE_COMMIT_RETAIN=5000"
Environment="STATE_BACKUP_DIR=/home/ubuntu/lumii-data/state-backups"
Environment="STATE_BACKUP_RETAIN=96"
Environment="STATE_BACKUP_MIN_INTERVAL_MS=300000"
Environment="LUMII_ADMIN_AUDIT_JOURNAL_PATH=/home/ubuntu/lumii-data/admin-audit-journal.jsonl"
Environment="LUMII_ADMIN_EXPORT_DIR=/home/ubuntu/lumii-data/admin-export-files"
ExecStartPre=/usr/bin/test -x /opt/node-v24.18.0-linux-x64/bin/node
ExecStartPre=/usr/bin/test -d /home/ubuntu/lumii-data
ExecStart=
ExecStart=/opt/node-v24.18.0-linux-x64/bin/node /home/ubuntu/lumii-project/scripts/lumii-backend.cjs
UMask=0077
